# dot_config/systemd/user/dotfiles-sync.timer
[Unit]
Description=Run guarded dotfiles synchronization every two hours

[Timer]
OnCalendar=*-*-* 00/2:00:00
Persistent=true
AccuracySec=1min
RandomizedDelaySec=3min
Unit=dotfiles-sync.service

[Install]
WantedBy=timers.target
