# dot_config/systemd/user/app-org.fcitx.Fcitx5@autostart.service.d/wayland-diagnose.conf
[Service]
# Keep GTK_IM_MODULE=fcitx for XWayland applications, but do not let Fcitx
# diagnose its own environment as an unnecessary Wayland override.
Environment=GTK_IM_MODULE=wayland
